# deploy/systemd/sw-data-sync.service
[Unit]
Description=Sync SW sewer dashboard data from the Raspberry Pi
Wants=network-online.target
After=network-online.target
OnFailure=sw-service-failure@%N.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/scripts/sync_sewer_data.py
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
