# deploy/systemd/sw-health-check.service
[Unit]
Description=Check SW sewer dashboard sync and backup freshness
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/scripts/check_sewer_health.py
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
