# HVAC/deploy/systemd/reverse-tunnel-hvac.service
[Unit]
Description=Reverse SSH tunnel for hvac.lucheestiy.com (VPS 18094 -> local 8094)
After=network-online.target
Wants=network-online.target

[Service]
User=root
Environment="AUTOSSH_GATETIME=0"
ExecStart=/usr/bin/autossh -M 0 -N \
  -o ExitOnForwardFailure=yes \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o StrictHostKeyChecking=accept-new \
  -i /root/.ssh/id_ed25519_inyp_vps \
  -R 127.0.0.1:18094:127.0.0.1:8094 root@104.236.97.60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
